# deploy/supervisor.conf
# Supervisor 配置 - FastAPI 进程管理
# 路径：/etc/supervisor/conf.d/langchain-agent.conf

[program:langchain-agent]
command=/opt/langchain-react-agent/venv/bin/python /opt/langchain-react-agent/start_api.py
directory=/opt/langchain-react-agent
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/langchain-agent.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/opt/langchain-react-agent/venv/bin"

# 优雅重启配置
stopwaitsecs=10
stopsignal=TERM
